# backend/infrastructure/smart-content-backend.service
[Unit]
Description=Smart Content Platform Backend (Node.js)
After=network.target

[Service]
Type=simple
User=shirka
WorkingDirectory=/home/shirka/public_html/backend
ExecStart=/usr/bin/node start.js
Restart=always
Environment=NODE_ENV=production
Environment=PORT=5001

[Install]
WantedBy=multi-user.target
